# docker/s6/migrations/up
/command/with-contenv sh -c "cd /app/NightSky; python manage.py migrate"
